# schedule_reports.cron
14 15 * 2,5,8,11 4,6 root /scripts/schedule_reports_filter.sh && /scripts/reports.sh
